# digital-lifeline/rootfs/etc/services.d/digital_lifeline/run
#!/bin/bash

TARGET="/config/custom_components/digital_lifeline"
SOURCE="/integration/custom_components/digital_lifeline"

echo "[INFO] Digital Lifeline: integratie installeren..."

if [ -d "$TARGET" ]; then
    echo "[INFO] Bestaande installatie gevonden – bijwerken..."
    rm -rf "$TARGET"
fi

mkdir -p "$TARGET"
cp -r "$SOURCE/." "$TARGET/"

echo "[INFO] -----------------------------------------------------"
echo "[INFO] Installatie geslaagd!"
echo "[INFO]"
echo "[INFO] Volgende stap:"
echo "[INFO]   1. Herstart Home Assistant"
echo "[INFO]   2. Ga naar Instellingen → Apparaten & Diensten"
echo "[INFO]   3. Klik op 'Integratie toevoegen'"
echo "[INFO]   4. Zoek op 'Digital Lifeline' en volg de stappen"
echo "[INFO] -----------------------------------------------------"

while true; do sleep 3600; done
