# digital-lifeline/rootfs/etc/services.d/digital_lifeline/run
#!/usr/bin/with-contenv bashio

TARGET="/config/custom_components/digital_lifeline"
SOURCE="/integration/custom_components/digital_lifeline"

bashio::log.info "Digital Lifeline: integratie installeren..."

if [ -d "$TARGET" ]; then
    bashio::log.info "Bestaande installatie gevonden – bijwerken..."
    rm -rf "$TARGET"
fi

mkdir -p "$TARGET"
cp -r "$SOURCE/." "$TARGET/"

bashio::log.info "-----------------------------------------------------"
bashio::log.info "Installatie geslaagd!"
bashio::log.info ""
bashio::log.info "Volgende stap:"
bashio::log.info "  1. Herstart Home Assistant"
bashio::log.info "  2. Ga naar Instellingen → Apparaten & Diensten"
bashio::log.info "  3. Klik op 'Integratie toevoegen'"
bashio::log.info "  4. Zoek op 'Digital Lifeline' en volg de stappen"
bashio::log.info "-----------------------------------------------------"

while true; do
    sleep 3600
done
